# workers/crontab
* * * * * /Users/student/.nvm/versions/node/v6.9.5/bin/node /Users/student/code/hrsf81-web-historian/workers/htmlfetcher.js
